# deployment/sysmon-agent.service
# =========================================================================
#  sysmon-agent — systemd service unit (Native Bash Agent)
# =========================================================================
#
# Install to: /lib/systemd/system/sysmon-agent.service
# Then run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now sysmon-agent.service
#
# =========================================================================

[Unit]
Description=Sysmon Lightweight Monitoring Agent (Native Bash)
Documentation=https://github.com/Vuvuive22/Lightweight-System-Monitoring-Agent
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple

ExecStart=/bin/bash /opt/sysmon-agent/agent.sh
WorkingDirectory=/opt/sysmon-agent

Restart=always
RestartSec=5

# Logging — stdout/stderr go to the systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sysmon-agent

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
